# dot-config/systemd/user/hypridle.service.d/log.conf
[Service]
StandardOutput=append:%h/.local/state/hypridle.log
StandardError=inherit
