# dot-config/systemd/user/hypridle.service.d/log.conf
[Service]
ExecStart=
ExecStart=/bin/sh -c "exec /usr/bin/hypridle 2>&1 | gawk '{ print strftime(\"%%FT%%T%%z\"), $$0; fflush() }'"
StandardOutput=append:%h/.local/state/hypridle.log
StandardError=inherit
