# dot-config/systemd/user/pi-dream.timer
[Unit]
Description=Run Pi Dream nightly

[Timer]
OnCalendar=*-*-* 03:30:00
RandomizedDelaySec=45m
Persistent=true

[Install]
WantedBy=timers.target
